# aw_sync/lite.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd template for the standalone aw_sync/lite.py script. install_lite.sh
  substitutes __HOME__ → $HOME and __SCRIPT__ → the script's full path,
  then writes the result to ~/Library/LaunchAgents/io.lifeos.awsync.plist.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.lifeos.awsync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__SCRIPT__</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/aw_sync.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/aw_sync.log</string>
</dict>
</plist>
