# tools/com.lifeos.calai-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent: nightly Cal AI nutrition sync from the iPhone (fully unattended).

  Backs up the phone with libimobiledevice into a non-TCC dir (no Full Disk
  Access needed), extracts Cal AI's Model.sqlite, and ingests the diary
  (idempotent). The first backup is a full mirror; nightly runs are incremental.

  Install:
    cp tools/com.lifeos.calai-sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.lifeos.calai-sync.plist
  Logs: ~/Library/Logs/calai_sync*.log

  Requires: brew install libimobiledevice; the iPhone paired + reachable
  (USB or Wi-Fi sync) and unlocked at run time; backup encryption OFF.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lifeos.calai-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/santiagoaday/Documents/LifeOS/tools/calai_sync.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CALAI_USE_IDEVICEBACKUP</key><string>1</string>
    <key>CALAI_BACKUP_ROOT</key><string>/Users/santiagoaday/CalAIBackups</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- Poll hourly: iOS only backs up an UNLOCKED phone, so we can't pick a fixed
       time. The script skips quietly when the phone is absent/locked and syncs
       (at most once/18h) the first time it's connected + unlocked. -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>/Users/santiagoaday/Library/Logs/calai_sync.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/santiagoaday/Library/Logs/calai_sync.err.log</string>
</dict>
</plist>
